# gestures/gesture-daemon.conf
start on started unity8
console log
exec /home/phablet/bin/gesture-daemon.py 2>&1 | logger -t gesture-daemon
